# download_storage.mount
[Mount]
What=/dev/nvme1n1
Where=/var/www/download.gluster.org/
