# mailman_storage.mount
[Unit]
Conflicts=umount.target
Before=local-fs.target umount.target
After=swap.target

[Mount]
What=/dev/nvme2n1
Where=/var/www/lists.gluster.org/
